# supervisor.conf
[program:website]
command = /usr/local/sbin/start.sh
[inet_http_server]
port = :9001
[program:gotty]
command = gotty -w sh